# website/backend/album_of_the_day_backend/album_of_the_day_backend/website/tasks/systemd/update_genre_descriptions.service
[Unit]
Description=Updates genre descriptions from online sources using the task runner.
Requires=update_genre_descriptions.timer
[Service]
Type=oneshot
ExecStart=/usr/bin/task_runner update_genre_descriptions
[Install]
WantedBy=multi-user.target