# website/backend/album_of_the_day_backend/album_of_the_day_backend/website/tasks/systemd/update_genre_descriptions.timer
[Unit]
Requires=update_genre_descriptions.service
[Timer]
Unit=update_genre_descriptions.service
OnCalendar=*-*-* 03:00
[Install]
WantedBy=timers.target
