# website/backend/album_of_the_day_backend/album_of_the_day_backend/website/tasks/systemd/update_daily_rotations.timer
[Unit]
Requires=update_daily_rotations.service
[Timer]
Unit=update_daily_rotations.service
OnCalendar=*-*-* 23:55:00
[Install]
WantedBy=timers.target