# website/backend/album_of_the_day_backend/album_of_the_day_backend/website/tasks/systemd/update_spotify_uids.service
[Unit]
Description=Updates Spotify UIDs.
Requires=update_spotify_uids.timer
[Service]
ExecStart=/usr/bin/task_runner update_spotify_uids
Type=oneshot
[Install]
RequiredBy=multi-user.target