# website/backend/album_of_the_day_backend/album_of_the_day_backend/website/tasks/systemd/update_spotify_uids.timer
[Unit]
Requires=update_spotify_uids.service
[Timer]
Unit=update_spotify_uids
OnCalendar=*-*-* 03:15
[Install]
WantedBy=timers.target
