# website/backend/album_of_the_day_backend/album_of_the_day_backend/website/tasks/systemd/update_album_covers.timer
[Unit]
Description=Updates album covers from Last.FM every day at 03:30.
Requires=update_album_covers.service
[Timer]
Unit=update_album_covers
OnCalendar=*-*-* 03:30
[Install]
WantedBy=timers.target
